# lvm-init.conf
description "Handle AWS instance storage"

start on filesystem
task

script
    /sbin/lvm-init
end script
